# scripts/com.luca.llm-protocol-gateway.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.luca.llm-protocol-gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__ROOT__/scripts/run-gateway-service.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>GATEWAY_ADDR</key>
    <string>0.0.0.0:18093</string>
    <key>PATH</key>
    <string>__HOME__/.bun/bin:__HOME__/.local/bin:__HOME__/.opencode/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__ROOT__/.cache/gateway-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__ROOT__/.cache/gateway-launchd.err.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
